# test.conf
[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[supervisord]
nodaemon=false
user=root
http_port=/var/run/supervisor.sock
pidfile=/var/run/micro/microca.pid
logfile=/var/log/microca/microca.log

[program:myapp]
command=/opt/microCA/venv/bin/gunicorn app:app -b 0.0.0.0:8898 -w 4
directory=/opt/microCA/venv/bin
startsecs=0
stopwaitsecs=0
autostart=false
autorestart=false
stdout_logfile=/opt/microCA/logs/web.log 
stderr_logfile=/opt/microCA/logs/web.err
